# crontab
*/5 * * * * django-user /app/.venv/bin/python /app/telethon.py >> /var/log/cron.log 2>&1
